# root/etc/services.d/kettle/run
#!/usr/bin/with-contenv sh

cd /kettle
exec python -u appkettle_mqtt.py --mqtt $MQTT_BROKER $MQTT_PORT $KETTLE_IP $KETTLE_IMEI;